# infrastructure/cron/crontab
# Run PHP scripts
* * * * * echo "[$(date)] Starting cron job" >> /var/log/cron.log 2>&1

* * * * * /usr/local/bin/php /var/www/html/bin/console app:healthcheck:logs >> /var/log/cron.log 2>&1
